# systemd/todolist_clean.service
[Unit]
Description=Todo List Cleanup Script
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/todo-list/scripts/clean.sh

[Install]
WantedBy=multi-user.target

