# systemd/todolist_webserver.service
[Unit]
Description=Todo List Server

[Service]
Type=simple
ExecStart=/opt/todo-list/run.sh
Restart=always

[Install]
WantedBy=multi-user.target

