# systemd/todolist_clean.timer
[Unit]
Description=Periodic Timer for HLS Heatbeat Check

[Timer]
OnBootSec=15min
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target

